# debian/ke-recv.upstart
description "starting ke-recv"
author "Philippe De Swert"
start on startup
stop on starting shutdown
console none
script
        export RUNLEVEL=2
        exec /etc/init.d/ke-recv start
end script
