# debian/ke-recv.upstart
description "starting ke-recv"
author "Philippe De Swert"

stop on ACT_DEAD
stop on TEST
stop on stopping hal
start on started hal
console none

script
        export RUNLEVEL=2
        exec /etc/init.d/ke-recv start
end script
